# systemd/journal-linker-voice-retry.service
# Transient voice failures — same role as launchd VoiceRetry.
[Unit]
Description=Journal Linker — voice retry (transient failures)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/journal-linker/journal-linker.env
ExecStart=/bin/bash ${JOURNAL_LINKER_REPO}/scripts/voice_retry.sh
