# systemd/journal-linker-feedback-sender.service
[Unit]
Description=Journal Linker — Telegram intent feedback long-polling daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/journal-linker/env
ExecStart=/bin/bash ${JOURNAL_LINKER_REPO}/scripts/feedback_sender.sh --daemon
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
