# systemd/journal-linker-feedback-sender.timer
[Unit]
Description=Journal Linker — intent feedback sender every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Unit=journal-linker-feedback-sender.service

[Install]
WantedBy=timers.target
