# launchd/DailyReflection.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>SCRIBE_JOURNAL_DIR</key>
		<string>REPLACE_WITH_YOUR_JOURNAL_PATH</string>
		<key>SCRIBE_MODEL</key>
		<string>llama3.1:8b</string>
		<key>SCRIBE_CTX</key>
		<string>8192</string>
		<key>SCRIBE_PUSHOVER_APP_TOKEN</key>
		<string>REPLACE_WITH_YOUR_PUSHOVER_APP_TOKEN</string>
		<key>SCRIBE_PUSHOVER_USER_KEY</key>
		<string>REPLACE_WITH_YOUR_PUSHOVER_USER_KEY</string>
		<!-- Optional -->
		<!-- <key>SCRIBE_PUSHOVER_DEVICE</key><string>iphone</string> -->
		<!-- <key>SCRIBE_PUSHOVER_SERVER</key><string>https://api.pushover.net</string> -->
		<!-- <key>SCRIBE_PUSHOVER_PRIORITY</key><string>0</string> -->
		<!-- <key>SCRIBE_DAILY_REFLECTION_WINDOW_START</key><string>16:00</string> -->
		<!-- <key>SCRIBE_DAILY_REFLECTION_WINDOW_END</key><string>21:00</string> -->
	</dict>
	<key>Label</key>
	<string>com.journal-linker.daily-reflection</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>REPLACE_WITH_REPO_PATH/scripts/daily_reflection.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>900</integer>
	<key>StandardErrorPath</key>
	<string>/tmp/daily-reflection.launchd.wrapper.err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/daily-reflection.launchd.wrapper.out.log</string>
	<key>WorkingDirectory</key>
	<string>REPLACE_WITH_REPO_PATH</string>
</dict>
</plist>
